# flask-demo.sup.ini
[program:flask-demo-5002]
; for supervisor deploy
directory=/opt/deploy/rtsf/flask-demo
command=python3 manager.py run

numprocs=1
autostart=true
autorestart=true

stopsignal=INT
